# heartbeat.crontab
# m h dom mon dow command
* * * * * /etc/network/if-up.d/heartbeat
